# jack_connection@jack_out.front\x2dleft\x2cspeakers_bus_volume.in_1.service.d/override.conf
[Unit]

Wants = jack_client@pulseaudio.service
After = jack_client@pulseaudio.service
Wants = jack_volume@speakers_bus.service
After = jack_volume@speakers_bus.service
