# mpv_ipc@.service.d/sonarr.conf
[Service]
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=de5afec0d57a43c4bba52aad5bfae5cd
