# alsa_out@amd.service.d/override.conf
[Service]
Environment="ALSA_DEVICE=hdmi:CARD=HDMI,DEV=0"
