# jack_volume@jack_out.service.d/override.conf
[Unit]
Description=Volume control for PulseAudio default sink input

[Service]
Environment="DISPLAY_NAME=Linux"
Environment=MIXER_GROUP=input
