# hdmi_input@2.service.d/override.conf
[Unit]
Wants = vrroom_notify.service

[Service]
Environment="SOURCE=4090"
