# hdmi_input@2.service.d/override.conf
[Unit]
Wants = vrroom.target

[Service]
Environment="SOURCE=4090"
