# outlet@8.service.d/override.conf
[Unit]

Description = PDU outlet 8: Desk (always on, ExecStop disabled)
StopWhenUnneeded = no

[Service]

ExecStop=

[Install]

WantedBy = default.target
