# outlet@8.service.d/override.conf
[Unit]

Description = PDU outlet 8: Desk (always on, never turned off)
StopWhenUnneeded = no

[Service]

Environment=PDU_NO_OFF_ON_STOP=1

[Install]

WantedBy = default.target
