# outlet@8.service.d/override.conf
[Unit]

Description = Desk
StopWhenUnneeded = no

[Service]

ExecStop=

[Install]

WantedBy = default.target
