# outlet@8.service.d/override.conf
[Unit]

Description = Outlet 8
StopWhenUnneeded = yes

# Used to be ferrofish
