# ad4d.target
[Unit]
Description=Shure AD4D — Axient Digital wireless mic receiver Dante configuration
Documentation=file://%h/bin/dante_device

Wants=dante-config@ad4d.target

[Install]
WantedBy=multi-user.target
