# audio-control-panels.target

[Unit]

Description = Audio control panels

Wants = cadence.service
Wants = catia.service
Wants = pavucontrol.service
