# audio.slice
[Unit]
Description = Audio processing slice

[Slice]
IOWeight=200
