# avio-aes3-1.target
[Unit]
Description=Audinate AVIO AES3-1 — AES3-to-Dante adapter connected to HDFury VRRoom HDMI audio output
Documentation=file://%h/bin/dante_device

Wants=dante-config@avio-aes3-1.target

[Install]
WantedBy=audio-baseline.target
