# avio-bt-1.target
[Unit]
Description=Audinate AVIO BT-1 — Bluetooth-to-Dante adapter, receives phone/Switch/Steam Deck audio
Documentation=file://%h/bin/dante_device

Wants=dante-config@avio-bt-1.target

[Install]
WantedBy=audio-baseline.target
