# avio-usb-1.target
[Unit]
Description=Audinate AVIO USB-1 — USB-to-Dante adapter connected to Windows gaming PC
Documentation=file://%h/bin/dante_device

Wants=dante-config@avio-usb-1.target

[Install]
WantedBy=audio-baseline.target
