# b.service
[Unit]
Description=Service B
# Requires=a.service
After=a.service

[Service]
ExecStart=/usr/bin/sleep infinity
