# cable@satellite-stereo-1-L_portico-left-A-input.target
[Unit]
Description=XLR: Satellite 5059 STEREO 1 L → Portico 5043 Left A INPUT
DefaultDependencies=no
