# cable@a32-analog-output-1-8_satellite-line-inputs-1-8.target
[Unit]
Description=DSub25: A32 ANALOG OUTPUT 1-8 → Satellite 5059 LINE INPUTS 1-8
DefaultDependencies=no
