# cable@shelford-line-out_mamba-A1.target
[Unit]
Description=XLR: Shelford Channel LINE OUT (0dB) → Mamba patchbay A1
DefaultDependencies=no
