# cable@p10r-headphone_se846-R.target
[Unit]
Description=IEM cable: P10R+ HEADPHONE → SE846 R
DefaultDependencies=no
