# cable@portico-left-A-output_p10t-TX1-L-CH1.target
[Unit]
Description=XLR: Portico 5043 Left A OUTPUT → P10T TX1 L/CH1
DefaultDependencies=no
