# cable@eaton-ups-output-1_pdu-ac.target
[Unit]
Description=AC: Eaton UPS OUTPUT 1 → APC PDU
DefaultDependencies=no
