# cable@shelford-6db-out_mamba-A2.target
[Unit]
Description=XLR: Shelford Channel -6dB LINE OUT → Mamba patchbay A2
DefaultDependencies=no
