# cable@sm7b_shelford-mic-in.target
[Unit]
Description=XLR: SM7B → Shelford Channel mic input
DefaultDependencies=no
