# cable@wl185_ad1-1.target
[Unit]
Description=Lav cable: WL185 → AD1 bodypack
DefaultDependencies=no
