# cable@mamba-B1-8_a32-analog-input-1-8.target
[Unit]
Description=DSub25: Mamba patchbay B1-8 → A32 ANALOG INPUT 1-8
DefaultDependencies=no
