# cable@satellite-stereo-1-R_portico-left-B-input.target
[Unit]
Description=XLR: Satellite 5059 STEREO 1 R → Portico 5043 Left B INPUT
DefaultDependencies=no
