# cable@eaton-ups-group-1-port-2_ugreen-charger-ac.target
[Unit]
Description=AC: Eaton UPS GROUP 1 PORT 2 → UGREEN USB charger
DefaultDependencies=no
