# cable@eaton-ups-group-2-port-1_shure-battery-charger-ac.target
[Unit]
Description=AC: Eaton UPS GROUP 2 PORT 1 → Shure battery charger
DefaultDependencies=no
