# cable@eaton-ups-output-2_www-server-ac.target
[Unit]
Description=AC: Eaton UPS OUTPUT 2 → www server
DefaultDependencies=no
