# cable@portico-left-B-output_p10t-TX1-R-CH2.target
[Unit]
Description=XLR: Portico 5043 Left B OUTPUT → P10T TX1 R/CH2
DefaultDependencies=no
