# chromium.slice
[Unit]
Description=Chromium browsers

[Slice]
MemoryMax=64G
MemoryHigh=60G
IOWeight=50
