# clipboard_sync.service
[Unit]

Description=Clipboard sync
Documentation=file://%h/bin/clipboard_sync
After=network.target

[Service]

Environment = DISPLAY=:0
Restart=always
RestartSec=10
Type=simple

ExecStart=%h/.local/bin/clipboard-sync --redis-host 127.0.0.1

[Install]

WantedBy=default.target
