# codex-conversations-recent.target
[Unit]
Description=Recent Codex conversations
StopWhenUnneeded=no
