# codex-conversations.target
[Unit]
Description=Codex conversations
StopWhenUnneeded=no

[Install]
WantedBy=xmonad.target
