# cpu-powersave.target
[Unit]
Description=CPU power saving mode
Conflicts=cpu-normal.target
PartOf=powersave.target
Wants=cpu-powersave.service
After=cpu-powersave.service
