# cpu-normal.target
[Unit]
Description=CPU normal performance mode
Conflicts=cpu-powersave.target
Wants=cpu-normal.service
After=cpu-normal.service
