# dante-config@.target
[Unit]
Description=Dante configuration enforcement — %i
PartOf=dante-enforcement.target
