# dante-device-sample-rate@.target
[Unit]
Description=Dante sample rate: %I
Documentation=file://%h/bin/dante_device
Wants=dante-manager.service
PartOf=dante-manager.service
