# dante-device-name@.target
[Unit]
Description=Dante device name: %I
Documentation=file://%h/bin/dante_device
Wants=dante-manager.service
PartOf=dante-manager.service
