# discord.target
[Unit]

Description=Discord — app with audio routing
StopWhenUnneeded = no

Wants = discord.service
Wants = discord-audio.target
