# discord.service
[Unit]

Description = Discord
Documentation=https://discord.com
Requires = sound.target graphical-session.target
Wants = jack_client@discord.service
Wants = discord_in.target
After = jack_client@discord.service
BindsTo = discord.target
# StopWhenUnneeded = yes

[Service]
Nice = 19
IOSchedulingClass=idle

# Limit Discord memory usage to prevent system freezes
MemoryMax=8G
MemoryHigh=6G

# Redirect logs to volatile tmpfs
StandardOutput=file:/run/user/%U/discord.log
StandardError=file:/run/user/%U/discord.log

Environment = DISPLAY=:0
ExecStart = /usr/bin/discord
ExecStop = /usr/bin/killall -I -9 discord

KillMode=process
Restart=no
RestartSec=6
TimeoutStartSec=2s
TimeoutStopSec=2s

[Install]
WantedBy = default.target
