# discord.target
[Unit]

Description=Discord
StopWhenUnneeded = no

Wants = discord.service
