# discord.target
[Unit]

Description=Discord
StopWhenUnneeded = no

Wants = discord.service
Wants = discord_in.target
Wants = discord_out.target
