# display-monitors.target

[Unit]

PropagatesReloadTo = display-monitors.service
PropagatesReloadTo = xmonad.target
Wants = display-monitors.service

