# eaton-ups.target
[Unit]
Description=Eaton UPS — lithium iron 1350W uninterruptible power supply
DefaultDependencies=no
Requires=eaton-ups-monitor.service
After=eaton-ups-monitor.service
BindsTo=eaton-ups-monitor.service
