# eaton-ups-monitor.service
[Unit]
Description=Eaton UPS Power Monitor — shared SNMP collector and trap handling
Documentation=file://%h/bin/eaton_ups_snmp
After=network-online.target

[Service]
ExecStart=%h/.local/bin/surface power ups monitor
Restart=always
RestartSec=5
Type=notify
Environment=UPS_HOST=192.168.1.6
Environment=CAPACITY_THRESHOLD=95
Environment=POLL_INTERVAL=2
Environment=SURFACE_REDIS_HOST=192.168.1.107
Environment=SURFACE_REDIS_DB=3
Environment=ALERT_EMAIL=s00pcan@gmail.com
Environment=COLLECTOR_LEASE_SECONDS=15
Environment=COLLECTOR_CHECK_INTERVAL=5

[Install]
WantedBy=default.target
