# ecasound@.service
[Unit]

Description=Ecasound effects processor — LADSPA compressor/effects chain for %I
Documentation=man:ecasound(1)
StopWhenUnneeded=yes
After=jackd.service
BindsTo=jackd.service
AssertFileIsExecutable=/usr/bin/ecasound

[Service]

Restart=always
RestartMode=direct
RestartSec=5
TimeoutStartSec=50s
TimeoutStopSec=10s

Slice=audio-rt.slice

CPUAffinity=1-4
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=90
CPUSchedulingResetOnFork=no

LimitCORE=infinity
LimitMEMLOCK=infinity
LimitRTPRIO=90
LimitRTTIME=infinity

[Install]
WantedBy=default.target
