# ferrofish-a32.target
[Unit]

Description=Ferrofish A32 — Dante/ADAT/analog audio interface
Documentation=file://%h/bin/dante_device

After=ferrofish-a32.service
BindsTo=ferrofish-a32.service
Requires=ferrofish-a32.service

Wants=dante-config@a32.target
