# ferrofish-a32.service
[Unit]
Description=Ferrofish A32 — Dante/ADAT/analog audio interface
Documentation=file://%h/bin/dante_device
Requires=outlet@4.service
After=outlet@4.service
BindsTo=outlet@4.service
PartOf=ferrofish-a32.target
StopWhenUnneeded=no
RefuseManualStart=no
RefuseManualStop=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/zsh -ic "sleep 4.5"

ExecStop=/usr/bin/zsh -ic "echo 'Shutting down Ferrofish A32...'"

TimeoutStartSec=25s
TimeoutStopSec=10s
RestartSec=10s
