# ferrofish-a32.target

[Unit]

BindsTo = outlet@3.service
Description = Ferrofish A32 Dante
PartOf = wireless-audio@1.service
PartOf = wireless-audio@2.service
Requires = outlet@3.service
