# ferrofish-a32.target
[Unit]

Description=Ferrofish A32 Dante

# After=dante-device@a32.service
After= ferrofish-a32.service
# BindsTo=dante-device@a32.service
BindsTo=ferrofish-a32.service
PartOf=wireless-audio@1.service
PartOf=wireless-audio@2.service
# Requires=dante-device@a32.service
Requires=ferrofish-a32.service
