# game-affinity.service
[Unit]
Description=Game CPU Affinity Daemon — pins game processes to isolated high-clock cores
Documentation=file://%h/bin/game-affinity

[Service]
Type=simple
ExecStart=%h/.local/bin/game-affinity
Restart=on-failure

[Install]
WantedBy=default.target
