# gaming.target
[Unit]
Description=Gaming mode — 8BitDo Pro 2 controller connected
Wants=8bitdo-pro2.target
StopWhenUnneeded=no
