# gaming.target
[Unit]
Description=Gaming Mode Active
Wants=8bitdo-pro2.target
StopWhenUnneeded=no
