# hdd-monitor.service
[Unit]
Description=HDD power state and temperature monitor
Documentation=file://%h/bin/surface

[Service]
Type=notify
ExecStart=%h/.local/bin/surface-hdd-monitor
Restart=always
Environment=POLL_INTERVAL=60

[Install]
WantedBy=default.target
