# hdmi-4k120.service
[Unit]
Description=Deprecated HDMI-0 4K120 output enforcement; use display-4k120@HDMI-0.service

[Service]
Type=notify
Environment=DISPLAY=:0
ExecStart=%h/.local/bin/surface-hdmi-enforce HDMI-0
Restart=always
RestartSec=5
