# hdmi-4k120.service
[Unit]
Description=Enforce HDMI 4K 120Hz display mode
Documentation=file://%h/bin/hdmi-4k120
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/bin/hdmi-4k120
Environment=DISPLAY=:0
