# hdmi-4k120.service
[Unit]
Description=Surface HDMI 4K120 output enforcement
ConditionEnvironment=DISPLAY

[Service]
Type=notify
Environment=DISPLAY=:0
ExecStart=%h/.local/bin/surface-hdmi-enforce
Restart=always
RestartSec=5
