# hdmi-4k120.timer
[Unit]
Description=Enforce HDMI 4K120 periodically

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
