# hhkb.service
[Unit]

Description = Happy Hacking Keyboard setup
Documentation=man:xset(1)
StopWhenUnneeded = no

[Service]

ExecStart = /usr/bin/sh -c '/usr/bin/xset r rate 200 125 || true'
ExecReload = /usr/bin/sh -c '/usr/bin/xset r rate 200 125 || true'
RemainAfterExit = no
Restart = no
RestartSec = 10m
StandardOutput = null
TimeoutStopSec = 1
Type = simple

[Install]

WantedBy = default.target
