# hhkb.service

[Unit]
Description=Happy Hacking Keyboard

[Service]
Environment=DISPLAY=:0
RemainAfterExit=yes
ExecStart=/usr/lib/udev/hhkb.zsh
Restart=always
RestartSec=1
ExecStop=
Type=forking
User=chris
WorkingDirectory=/home/chris

[Install]
WantedBy=multi-user.target

