# hhkb.service

[Unit]
Description=Happy Hacking Keyboard

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/xset r rate 200 120
RestartSec=2
Type=oneshot

[Install]
WantedBy=graphical.target

