# infrastructure-map.target
[Unit]
Description = Infrastructure Map
Wants = infrastructure-map.service
