# jack_clients.target

[Unit]

After = jackd.service
Description = Jack clients
