# jackd.target

[Unit]

Description = JACK Audio Connection Kit sound server
PropagatesReloadTo = jackd.service
Wants = jackd.service
