# jack_effects.target
[Unit]

Description = Watch jack over dbus to connect to the speakers ecasound effects
StopWhenUnneeded = yes

Wants = jack_out_volume.target
