# jack_clients.target

[Unit]

After = jackd.service
Description = Jack clients
PartOf = jack.target
Requires = jackd.service

[Install]

WantedBy = jack.target
