# jack_client@tv_speakers.service
[Unit]

StopWhenUnneeded = yes
DefaultDependencies = no
Description = TV speakers — HDMI audio output
Documentation=man:alsa_out(1)

After = jackd.service
BindsTo = jackd.service

# Wants = tv_speakers_notify.service
Wants = ecasound@speakers.service
Wants = tv_speakers.target

[Service]

Restart = always

ExecStart = /usr/bin/alsa_out -p 1024 -n 3-c2 -j '%i' -d 'hdmi:CARD=NVidia,DEV=0'
# ExecStart = /usr/bin/alsa_out -c2 -j '%i' -d 'hdmi:CARD=HDMI,DEV=0'

# speaker-test -D hw:2,3 -c2

[Install]

WantedBy=jack_clients.target
