# jack_in.target
[Unit]
Description=PulseAudio default source routing — jack_in capture path for generic Linux apps
After = jackd.service
BindsTo = jackd.service
