# jack_client@dante_pcie_in.service

[Unit]

# After = sound.target basic.target microphone.service
After = dante_pcie_in.service
Description = Dante PCIe capture — LX-Dante to JACK input bridge
Documentation=man:alsa_in(1)
BindsTo = dante_pcie_in.service
Requires = sound.target basic.target
# Requires = sound.target basic.target jackd.service
StopWhenUnneeded = no

[Service]
Type = notify

ExecStart = %h/.local/bin/surface bridge jack dante capture %i --channels 16 --rate 48000 --period 1024 --nperiods 2 --device asym_dante0

Restart = always
RestartSec = 5
TimeoutStartSec = 50

[Install]

WantedBy = jack_clients.target
