# jack-mixer-web.service
[Unit]
Description=JACK Volume Mixer Web UI — browser-based volume controls for all JACK volume processors
Documentation=file://%h/bin/jack_mixer_web
After=jackd.service
BindsTo=jackd.service
StopWhenUnneeded=yes

[Service]
Type=simple
ExecStart=%h/.local/bin/jack-mixer-web
Restart=always
RestartSec=3
