# jack_in.target
[Unit]
Description=Watch JACK over DBus to connect the default pulseaudio jack_in
After = jackd.service
BindsTo = jackd.service
