# jack_clients.target
[Unit]

After = jackd.service
Description = Jack clients
BindsTo = jack.target

[Install]

WantedBy = jack.target
