# jackdbus.service
[Unit]
Description=Mask JACK D-Bus daemon to prevent auto-start
Documentation=man:jackdbus(1)
After=dbus.service
Requires=dbus.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=dbus
BusName=org.jackaudio.service
ExecStart=/usr/bin/jackdbus auto
Restart=on-failure
RestartSec=5

# Capture all output to systemd journal
StandardOutput=journal
StandardError=journal

# Realtime scheduling
LimitRTPRIO=infinity
LimitRTTIME=infinity
LimitMEMLOCK=infinity

[Install]
WantedBy=default.target
