# jack.target

[Unit]

Description = Jack sound service
Requires = jackd.service 

[Install]

WantedBy = default.target

