# keyboard-device.path
[Unit]
Description = Path for HHKB

[Path]

PathChanged = /dev/input/hhkb
Unit = setup_keyboard.service

[Install]
WantedBy=default.target
