# macbook_personal_direct.target
[Unit]
Description=MacBook Personal direct audio routing (system:capture_11/12 → jack_effects)
After=jackd.service
BindsTo=jackd.service
Conflicts=macbook_personal_volume.target

# Direct routing: system:capture_11/12 → jack_effects (music:left/right aliases)
Requires=jack_connection@system.capture_11\x2cjack_effects.in_1.service
Requires=jack_connection@system.capture_12\x2cjack_effects.in_2.service
After=jack_connection@system.capture_11\x2cjack_effects.in_1.service
After=jack_connection@system.capture_12\x2cjack_effects.in_2.service