# mamba-patchbay.target
[Unit]
Description=Mamba XLR Patchbay — 16-ch Neutrik combo passive patchbay
DefaultDependencies=no
