# mednafen@.service

[Unit]

Description = Run game %i

[Service]

Type = simple

# [Install]
#
# WantedBy=default.target
