# mednafen@.service

[Unit]

Description = Mednafen emulator — run %i

[Service]

Type = simple

# [Install]
#
# WantedBy=default.target
