# microphone.target
[Unit]

After = outlet@7.service
BindsTo = outlet@7.service
Description=DANTE Microphone
PropagatesReloadTo = microphone.service
Requires = microphone.service
BindsTo = microphone.service
