# microphone.target
[Unit]
Description = Microphone — Shure SM7B through Neve Shelford Channel preamp, Shure P10T IEM transmitter ready

After = shelford-channel.target
After = shure-p10t.target
Requires = shelford-channel.target
BindsTo = shelford-channel.target
