# microphone.target
[Unit]
Description = Microphone

After = shelford-channel.target
After = shure-p10t.target
PropagatesReloadTo = microphone.service
Requires = shelford-channel.target
BindsTo = shelford-channel.target
