# microphone.target
[Unit]

After = shelford-channel.target
BindsTo = shelford-channel.target
Description=Microphone
PropagatesReloadTo = microphone.service
Requires = microphone.service
BindsTo = microphone.service
