# monitor_btc_price.timer
[Unit]
Description=Monitor BTC-USD exchange rate

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
