# movie.service

[Unit]
Conflicts=shutdown.target
Description=watch a movie

; TODO: Use pid to stop daemon

[Service]
Environment=DISPLAY=:0
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/home/chris/.scripts/video_player/player.zsh
Restart=always
RestartSec=1
Type=simple
WorkingDirectory=/home/chris/.scripts/video_player/

[Install]
WantedBy=graphical.target

