# music.target
[Unit]
Description=Music — MPD direct to dedicated Dante output, no volume or effects
After=jackd.service
BindsTo=jackd.service
Conflicts=music_direct.target
Conflicts=music_volume.target

Requires=mpd@music.service
After=mpd@music.service
Requires=jack_disconnection@music.left\x2cmusic_effects.in_1.service
Requires=jack_disconnection@music.right\x2cmusic_effects.in_2.service
Requires=jack_disconnection@music.left\x2cmusic_volume.in_1.service
Requires=jack_disconnection@music.right\x2cmusic_volume.in_2.service
Requires=jack_disconnection@music_volume.out_1\x2cmusic_effects.in_1.service
Requires=jack_disconnection@music_volume.out_2\x2cmusic_effects.in_2.service

After=jack_disconnection@music.left\x2cmusic_effects.in_1.service
After=jack_disconnection@music.right\x2cmusic_effects.in_2.service
After=jack_disconnection@music.left\x2cmusic_volume.in_1.service
After=jack_disconnection@music.right\x2cmusic_volume.in_2.service
After=jack_disconnection@music_volume.out_1\x2cmusic_effects.in_1.service
After=jack_disconnection@music_volume.out_2\x2cmusic_effects.in_2.service

# Direct routing: music → system playback (bypassing all effects)
Requires=jack_connection@music.left\x2csystem.playback_15.service
Requires=jack_connection@music.right\x2csystem.playback_16.service
After=jack_connection@music.left\x2csystem.playback_15.service
After=jack_connection@music.right\x2csystem.playback_16.service
