# netaudio.service
[Unit]
Description=Netaudio Dante device discovery daemon
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/netaudio --log-level INFO server daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
