# netaudio.service
[Unit]
Description=Netaudio — Dante device discovery and monitoring daemon
Documentation=file://%h/bin/netaudio
After=network.target
PartOf=dante-enforcement.target

[Service]
Type=notify
ExecStart=%h/.local/bin/netaudio server start
Restart=on-failure
RestartSec=5

[Install]
