# netaudio-enforce@.service
[Unit]
Description=Netaudio enforcement: %i
After=netaudio.service
Requires=netaudio.service
PartOf=netaudio.service

[Service]
Type=notify
ExecStart=%h/.local/bin/netaudio-enforce dante %I
NotifyAccess=all
TimeoutStartSec=15
TimeoutStopSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=netaudio-enforce-%I

[Install]
WantedBy=default.target
