# network-map.target
[Unit]
Description = Network Map
Wants = network-map.service
