# night-recording-start.timer
[Unit]
Description=Start night recording at 23:30

[Timer]
OnCalendar=*-*-* 23:30:00
Persistent=true

[Install]
WantedBy=timers.target
