# night-recording-stop.service
[Unit]
Description=Stop night recording

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user stop night-recording.target
