# outlet-enforcer@.service
[Unit]
Description=Enforce outlet %i state via SNMP
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/uv run --project /home/chris/bin/pdu_snmp pdu-systemd-enforcer %i
