# outlet@.service
[Unit]

DefaultDependencies = no
Description = PDU outlet %i
Documentation=file://%h/bin/pdu_snmp
Wants = pdu-manager.service
After = pdu-manager.service

[Service]
TimeoutStartSec = 10
TimeoutStopSec = 10
RestartSec = 5
Restart = no
Type = notify

ExecStart=%h/.local/bin/pdu-outlet-command run %i outlet@%i.service

[Install]

WantedBy=default.target
