# plex-transcode-cpu.target
[Unit]
Description=Plex is transcoding — CPU encoding
BindsTo=plex-transcode-watch.service
