# portico-5043-right.target
[Unit]
Description=Rupert Neve Designs Portico 5043 Right — stereo compressor
Requires=portico-5043-right.service
After=portico-5043-right.service
BindsTo=portico-5043-right.service
