# portico-5043-right.target
[Unit]

After = outlet@7.service
BindsTo = outlet@7.service
Description = Rupert Neve Portico 5043 (Right)
PartOf = wireless-audio@2.service
