# resume.service
[Unit]
Description=Resume build service

[Service]
Type=oneshot
WorkingDirectory=%h/latex_resume
ExecStart=/usr/bin/pdflatex resume.tex
