# resume.path
[Unit]
Description=Watch for changes to resume.tex

[Path]
PathChanged=%h/latex_resume/resume.tex

[Install]
WantedBy=default.target
