# retroarch.service
[Unit]
Description=RetroArch — multi-system emulator frontend
Documentation=man:retroarch(6)

[Service]
Type=simple
ExecStart=/usr/bin/retroarch --verbose
Restart=on-failure

[Install]
WantedBy=default.target
