# se846.target
[Unit]
Description=Shure SE846 — in-ear monitors
DefaultDependencies=no
