# setup_jack_alias.service
[Unit]

Description=JACK port alias setup
Documentation=file://%h/bin/jack_alias

After=jackd.service jackdbus.service
BindsTo=jackd.service

[Service]
Type=notify
RemainAfterExit=yes
TimeoutStartSec=infinity
TimeoutStopSec=15
KillMode=mixed

ExecStart=%h/.local/bin/surface audio jack dante-alias --alias
ExecStop=%h/.local/bin/surface audio jack dante-alias --unalias
