# shelford-channel.target
[Unit]

After = outlet@7.service
BindsTo = outlet@7.service
Description = Shelford channel
Requires = outlet@7.service
BindsTo = microphone.target
