# shelford-channel.target
[Unit]

After = shelford-channel.service
BindsTo = shelford-channel.service
Description = Rupert Neve Designs Shelford Channel — mic preamp/EQ/compressor
Requires = shelford-channel.service
