# shure-p10t.target
[Unit]
Description = Shure P10T — wireless audio transmitter
After = ferrofish-a32.target
Requires = shure-p10t.service
After = shure-p10t.service
BindsTo=  shure-p10t.service
Before=wireless-audio@1.service
Before=wireless-audio@2.service
# BindsTo=wireless-audio@1.service
# BindsTo=wireless-audio@2.service
