# shure-p10r.target
[Unit]
Description=Shure P10R+ — wireless IEM bodypack receiver
DefaultDependencies=no
