# shure-battery-charger.target
[Unit]
Description=Shure SBC200 — dual docking recharging station
DefaultDependencies=no
