# shure-ad1-1.target
[Unit]
Description=Shure AD1 #1 — Axient Digital bodypack wireless transmitter
DefaultDependencies=no
