# sm7b.target
[Unit]
Description=Shure SM7B — dynamic microphone
DefaultDependencies=no
