# snmp_trap.service
[Unit]
Description=SNMP Trap Handler

[Service]
ExecStart=/usr/bin/python3 /home/chris/bin/pdu_snmp/src/pdu_snmp/trap.py
Restart=always
Type=notify
