# surface-tmux-workspace.target
[Unit]
Description=Restore saved tmux workspace
StopWhenUnneeded=no

[Install]
WantedBy=default.target
