# tv.target

[Unit]

Description = TV/media display

# PropagatesReloadTo = movie.service

After = xmonad.target
BindsTo = xmonad.target
Wants = movie.service
