# ugreen-charger.target
[Unit]
Description=UGREEN USB charger
DefaultDependencies=no
