# video.target
[Unit]
Description=Video/HDMI infrastructure

Wants=hdmi_input_detect.service
