# windows-sleep.target
[Unit]
Description=Windows gaming PC desired state: sleeping
Conflicts=windows-wake.target
Wants=windows-sleep.service
After=windows-sleep.service
