# windows-cpu-power-normal.service
[Unit]
Description=Windows gaming PC CPU power mode: normal
Documentation=man:ssh(1)
Conflicts=windows-cpu-power-powersave.service windows-cpu-power-blackout.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/usr/bin/zsh -c "until nc -vzw 1 192.168.1.194 22 2>/dev/null; do sleep 0.5; done"
ExecStart=/usr/bin/ssh windows-pwsh "powercfg /setacvalueindex 381b4222-f694-41f0-9685-ff5bb260df2e 54533251-82be-4824-96c1-47b60b740d00 75b0ae3f-bce0-45a7-8c89-c9611c25e101 0 && powercfg /setacvalueindex 381b4222-f694-41f0-9685-ff5bb260df2e 54533251-82be-4824-96c1-47b60b740d00 75b0ae3f-bce0-45a7-8c89-c9611c25e100 0 && powercfg /setacvalueindex 381b4222-f694-41f0-9685-ff5bb260df2e 54533251-82be-4824-96c1-47b60b740d00 be337238-0d82-4146-a960-4f3749d470c7 2 && powercfg /setactive a1841308-3541-4fab-bc81-f71556f20b4a && powercfg /setactive 381b4222-f694-41f0-9685-ff5bb260df2e"
ExecStartPost=/usr/bin/ssh windows-pwsh "powercfg /getactivescheme"
