# wireless-microphone.target
[Unit]

After = shure-ad4d.target
BindsTo = shure-ad4d.target
Description = Wireless microphone — Shure AD4D receiver
Requires = shure-ad4d.target
