# wireless-microphone.target

[Unit]

After = shure-ad4d.target
BindsTo = shure-ad4d.target
Description = Wireless microphone
Requires = shure-ad4d.target
