# wireless@p10t-tx1_p10r.target
[Unit]
Description=RF: P10T Transmitter 1 → P10R+ bodypack receiver
DefaultDependencies=no
