# wireless@ad1-1_ad4d.target
[Unit]
Description=RF: AD1 bodypack → AD4D receiver
DefaultDependencies=no
