# wl185.target
[Unit]
Description=Shure WL185 — lavalier condenser microphone (R183 omni cartridge)
DefaultDependencies=no
