# www-server.target
[Unit]
Description=www — primary server (i9-7940X)
DefaultDependencies=no
