# pulseaudio.service.d/override.conf
[Unit]
Wants=jack_client@pulseaudio.service
