# jack_volume@linux_main.service.d/override.conf
[Unit]
Description=Volume control for linux main mix bus

[Service]
Environment=CHANNELS=2
Environment="DISPLAY_NAME=Linux Main"
Environment=MIXER_GROUP=bus
