# wireless-audio@2.service.d/override.conf

[Unit]

Wants = outlet@5.service
After = outlet@5.service 

