# jack_connection@speakers_bus_volume.out_2\x2cspeakers_effects.in_2.service.d/override.conf
[Unit]
StopWhenUnneeded = yes

After = jack_volume@speakers_bus.service
After = ecasound@speakers.service
Wants = jack_volume@speakers_bus.service
Wants = ecasound@speakers.service
