# jack_volume@windows_gaming.service.d/override.conf
[Unit]
Description=Volume control for Windows gaming PC audio input

[Service]
Environment="DISPLAY_NAME=Windows AVIO 1"
Environment=MIXER_GROUP=input
