# outlet@5.service.d/override.conf

[Unit]

PartOf = wireless-audio@2.service
StopWhenUnneeded = yes
