# outlet@5.service.d/override.conf
[Unit]

Description = Rupert Neve Portico 5043 (Right) power outlet
BindsTo = portico-5053-right.target
StopWhenUnneeded = yes
