# outlet@5.service.d/override.conf
[Unit]

Description = PDU outlet 5: Neve Portico 5043 Left (stereo compressor)
StopWhenUnneeded = yes
