# outlet@5.service.d/override.conf
[Unit]

# Description = Rupert Neve Portico 5043 (Left) power outlet
# BindsTo = portico-5053-left.target

Description = Outlet 5
StopWhenUnneeded = yes
