# jack_volume@macbook_work.service.d/override.conf
[Service]
Environment="DISPLAY_NAME=MacBook Work"
