# jack_volume@macbook_work.service.d/override.conf
[Unit]
Description=Volume control for work MacBook audio input

[Service]
Environment="DISPLAY_NAME=MacBook Work"
Environment=MIXER_GROUP=input
