# jack_volume@macbook_personal.service.d/override.conf
[Unit]
Description=Volume control for personal MacBook audio input

[Service]
Environment="DISPLAY_NAME=MacBook Personal"
Environment=MIXER_GROUP=input
