# jack_volume@macbook_personal.service.d/override.conf
[Service]
Environment="DISPLAY_NAME=MacBook Personal"
