# jack_volume@music.service.d/override.conf
[Unit]
Description=Volume control for MPD music player input

[Service]
Environment=DISPLAY_NAME=Music
Environment=MIXER_GROUP=input
