# jack_volume@music.service.d/override.conf
[Service]
Environment=DISPLAY_NAME=Music
