# hdmi_input@1.service.d/override.conf
[Unit]

# Requires = vrroom.target

Conflicts = jack_disconnection@vrroom_volume.out_1\x2cjack_effects.in_1.service
Conflicts = jack_disconnection@vrroom_volume.out_2\x2cjack_effects.in_2.service

[Service]

# Environment="SOURCE=PS5"
Environment=SOURCE=
