# hdmi_input@1.service.d/override.conf
[Unit]
Description = HDMI input 1: unassigned (loose cable)

Conflicts = hdmi_input@0.service hdmi_input@2.service hdmi_input@3.service

Wants = vrroom_audio.target

[Service]

Environment=SOURCE=
