# hdmi_input@1.service.d/override.conf
[Unit]

# Requires = vrroom.target

[Service]

# Environment="SOURCE=PS5"
Environment=SOURCE=
