# volume_control@music.service.d/override.conf
[Unit]
Description=Volume control for Music (MPD) audio
After=setup_jack_alias.service

[Service]
Environment=CONTROL_PORT=7015
Environment=INSTANCE_NAME=music