# jack-connection-manager.service.d/health.conf
[Unit]
After=jackd.service jackdbus.service
Requires=jackd.service
StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
TimeoutStartSec=15
WatchdogSec=8
