# mednafen@super_mario_bros_3.service.d/override.conf

[Unit]

Description = Run Super Mario Bros 3 for NES
Wants = super_mario_bros_3_notify.service

[Service]
ExecStart = /media/backup/software/mednafen/mednafen/src/mednafen /media/games/nintendo/nes/games/super_mario_bros_3.nes
