# outlet@4.service.d/override.conf

[Unit]

StopWhenUnneeded = yes
# PartOf = wireless-audio@.service
# PartOf = wireless-audio@1.service
# PartOf = wireless-audio@2.service
