# outlet@4.service.d/override.conf
[Unit]

Description = PDU outlet 4: Ferrofish A32 (32-ch Dante audio interface)
StopWhenUnneeded = yes
