# outlet@4.service.d/override.conf
[Unit]

# Description = Rupert Neve Portico 5043 (Left)
# BindsTo = wireless-audio@1.service

Description = Outlet 4
StopWhenUnneeded = yes
