# volume_control@windows_gaming.service.d/override.conf
[Unit]
Description=Volume control for Windows Gaming audio
After=setup_jack_alias.service

[Service]
Environment=CONTROL_PORT=7007
Environment=INSTANCE_NAME=windows_gaming