# outlet@6.service.d/override.conf
[Unit]

Description = Rupert Neve Satellite 5059
PartOf = wireless-audio@1.service
PartOf = wireless-audio@2.service
StopWhenUnneeded = yes
