# wireless-audio@1.service.d/override.conf

[Unit]

Wants = outlet@6.service
After = outlet@6.service 
