# wireless-audio@1.service.d/override.conf
[Unit]

After = ferrofish-a32.target
After = microphone.target
After = portico-5043-left.target
Requires = ferrofish-a32.target
Requires = portico-5043-left.target
