# jack_volume@bluetooth.service.d/override.conf
[Service]
Environment=DISPLAY_NAME=Bluetooth
