# jack_volume@bluetooth.service.d/override.conf
[Unit]
Description=Volume control for Bluetooth audio input

[Service]
Environment=DISPLAY_NAME=Bluetooth
Environment=MIXER_GROUP=input
