# hdmi_input@0.service.d/override.conf
[Unit]

Conflicts = vrroom.target
# Wants = tv_speakers.target

[Service]

Environment=SOURCE=
