# hdmi_input@0.service.d/override.conf
[Unit]

Conflicts = vrroom_notify.service
# Wants = tv_speakers.target

[Service]

Environment=SOURCE=
