# jack_connection@jack_out.front\x2dleft\x2clinux_main_volume.in_1.service.d/override.conf
[Unit]

Wants = jack_client@pulseaudio.service
After = jack_client@pulseaudio.service
Wants = jack_volume@linux_main.service
After = jack_volume@linux_main.service
