# jack_disconnection@vrroom_volume.out_2\x2cspeakers_bus_volume.in_1.service.d/override.conf
[Unit]
Conflicts=jack_connection@vrroom_volume.out_2\x2cspeakers_bus_volume.in_1.service
