# hdmi_input@3.service.d/override.conf
[Unit]

Conflicts = hdmi_input@0.service hdmi_input@1.service hdmi_input@2.service

Wants = vrroom_audio.target

# [Service]
# Environment="SOURCE=Xbox"
