# hdmi_input@3.service.d/override.conf
[Unit]
Wants = vrroom_notify.service

[Service]
Environment="SOURCE=Xbox"
