# hdmi_input@3.service.d/override.conf
[Unit]
Wants = vrroom.target

[Service]
Environment="SOURCE=Xbox"
