# discord.service.d/sandbox.conf
[Service]
ProtectHome=tmpfs
BindPaths=%h/.config/discord
BindPaths=/run/user/%U/pulse
BindReadOnlyPaths=/run/user/%U/bus
BindReadOnlyPaths=%h/.config/pulse
BindReadOnlyPaths=%h/.Xauthority
BindReadOnlyPaths=%h/.local/share/fonts
BindReadOnlyPaths=%h/.config/fontconfig
BindReadOnlyPaths=%h/.icons
BindReadOnlyPaths=%h/.local/share/icons

ProtectSystem=full
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictRealtime=yes
LockPersonality=yes
SystemCallArchitectures=native
UMask=0077

NoNewPrivileges=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

CapabilityBoundingSet=
KeyringMode=private
ProtectProc=invisible
ProcSubset=pid

DevicePolicy=closed
DeviceAllow=/dev/dri/renderD128 rw
DeviceAllow=/dev/dri/card0 rw
DeviceAllow=/dev/snd rw
DeviceAllow=char-input r
