# ympd@books.service.d/override.conf

[Service]

Environment = YMPD_PORT=8096
