# outlet@2.service.d/override.conf
[Unit]
# Used to be P10T

Description = Shure AD4D
BindsTo = shure-ad4d.target
StopWhenUnneeded = yes
