# deprecated/jack_client@shure_sm7b.service

[Unit]
Description = Shure SM7B

[Service]
ExecStart=/usr/bin/alsa_in -c1 -j '%i' -d hw:Digital

[Install]
WantedBy=jack_clients.target
