# deprecated/mixed-microphones.service

[Unit]

Conflicts = microphone.service
Conflicts = speakers.service
# After = outlet@7.service
# After = outlet@2.service
Before = dante_in.service
# BindsTo = outlet@7.service
# BindsTo = outlet@2.service
DefaultDependencies = no
Description = Mixed microphones — dual-mic mode, conflicts with single-mic and speakers services
StopWhenUnneeded = yes
# RefuseManualStart = yes
# RefuseManualStop = yes
# Requires = dante_in.service
# Requires = outlet@7.service
# Requires = outlet@2.service

[Service]

ExecStart = /usr/bin/sleep 0.2
# ExecStartPost = /bin/zsh -ic "dante_mic_both off"
# ExecStartPost = /bin/zsh -ic "dante_mic_both off"
# ExecStartPost = /bin/zsh -ic "dante_mic_both on"
# ExecStartPost = /bin/zsh -ic "dante_mic_both on"
# ExecStartPost = /bin/zsh -ic "sleep 0.2 && dante_mic_both on"
# ExecStartPost = /bin/zsh -ic "sleep 0.2 && dante_mic_both on"

# ExecStop = /bin/zsh -ic "dante_mic_both off"

RemainAfterExit = yes
RestartSec=2
TimeoutStopSec=2s
# Type = oneshot
