# deprecated/movie.target
[Unit]

Description = Video player
Requires = movie.service
PropagatesReloadTo = movie.service

[Install]

WantedBy=default.target
