# mpv_ipc@slime.service.d/binding.conf
[Service]
Environment=SONARR_SERIES_ID=422
