# jack_volume@speakers_bus.service.d/override.conf
[Unit]
Description=Volume control for speakers mix bus

[Service]
Environment=CHANNELS=2
Environment="DISPLAY_NAME=Speakers"
Environment=MIXER_GROUP=bus
