# outlet@1.service.d/override.conf
[Unit]

Description = Outlet 1

[Service]

TimeoutStartSec = 40s
TimeoutStopSec = 15s
RestartSec = 5s
