# outlet@1.service.d/override.conf
[Unit]
# Used to be wall lamps

Description = Shure P10T
BindsTo = shure-p10t.target
StopWhenUnneeded = yes

[Service]

# ExecStartPost = /usr/bin/zsh -ic "sleep 4.5"
# ExecStartPost = /usr/bin/zsh -ic "until nc -vzw 1 'shure-p10t' 2202 2>/dev/null; do sleep 0.5; done"
# ExecStartPost = /usr/bin/zsh -ic '%h/.scripts/zsh/maybe_mute.zsh'

TimeoutStartSec = 40s
TimeoutStopSec = 5s
RestartSec = 5s
