# outlet@1.service.d/override.conf
[Unit]

Description = Outlet 1
StopWhenUnneeded = yes

[Service]

TimeoutStartSec = 40s
TimeoutStopSec = 5s
RestartSec = 5s
