# outlet@1.service.d/override.conf
[Unit]

Description = PDU outlet 1: Neve Satellite 5059 (16x2 summing mixer)

[Service]

TimeoutStartSec = 40s
TimeoutStopSec = 15s
RestartSec = 5s
