# jack_volume@discord_out.service.d/override.conf
[Unit]
Description=Volume control for Discord audio input

[Service]
Environment="DISPLAY_NAME=Discord Linux"
Environment=MIXER_GROUP=input
