# jack_volume@discord_out.service.d/override.conf
[Service]
Environment="DISPLAY_NAME=Discord Linux"
