# jack_volume@vrroom.service.d/override.conf
[Unit]
Description=Volume control for VRRoom HDMI audio input

[Service]
Environment=CHANNELS=2
Environment="DISPLAY_NAME=VRRoom HDMI"
Environment=MIXER_GROUP=input