# ympd@music.service.d/override.conf

[Service]

Environment = YMPD_PORT=8094
Environment = MPD_HOST=%t/mpd_music/socket
